# com.sean.moviedb.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.sean.moviedb</string>

    <key>ProgramArguments</key>
    <array>
      <!-- absolute paths, because PATH is nearly empty under launchd -->
      <string>/opt/homebrew/bin/npx</string>
      <string>serve</string>
      <string>-s</string>
      <string>/Users/sean/Sites/moviedb/dist/moviedb</string>
      <string>-l</string>
      <string>4200</string>
    </array>

    <!-- npx shells out to env(1) to find node; supply PATH so it can -->
    <key>EnvironmentVariables</key>
    <dict>
      <key>PATH</key>
      <string>/opt/homebrew/bin:/usr/bin:/bin</string>
    </dict>

    <!-- optional but handy -->
    <key>WorkingDirectory</key>
    <string>/Users/sean/Sites/moviedb/dist/moviedb</string>

    <!-- start at login & restart on crash -->
    <key>RunAtLoad</key><true/>
    <key>KeepAlive</key><true/>

    <!-- simple log files -->
    <key>StandardOutPath</key><string>/Users/sean/Library/Logs/moviedb.out</string>
    <key>StandardErrorPath</key><string>/Users/sean/Library/Logs/moviedb.err</string>
  </dict>
</plist>
